# bin/crontab
0 * * * * /projects/hashmweb/bin/update-script.sh >& /dev/null